# autostart/calised.desktop
[Desktop Entry]
Name=calised
TryExec=calised
Exec=calised
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
