# autostart/claws-mail.desktop
[Desktop Entry]
Name=claws-mail
TryExec=claws-mail
Exec=claws-mail
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
