# autostart/devilspie.desktop
[Desktop Entry]
Name=devilspie
TryExec=devilspie
Exec=devilspie
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
