# autostart/sogou.desktop
[Desktop Entry]
Name=sogou
TryExec=sogou-qimpanel
Exec=sogou-qimpanel
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
